# scrapy/crontab
0 5 * * 1-5 cd /code/erad; python run_scrapy_crawl_fund.py